# sample-daemon.service
[Unit]
Description=Simple daemon

[Service]
ExecStart=/usr/bin/java -jar /home/user/systemd-daemon-sample.jar
SuccessExitStatus=143
TimeoutStopSec=10